# launchd/com.yearflow.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yearflow.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/path/to/YearFlow/app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/YearFlow</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/path/to/YearFlow/logs/yearflow-launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/YearFlow/logs/yearflow-launchd.err.log</string>
</dict>
</plist>
